# containers/debian-client-slim/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
; Run the wrapper script which will, in turn, start sshd
command=/usr/local/bin/start-sshd-dynamic.sh
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

;[program:my-app]
; Replace this with the command for your other application
; command=/usr/local/bin/my-other-app
; autostart=true
; autorestart=true
; stdout_logfile=/dev/stdout
; stdout_logfile_maxbytes=0
; stderr_logfile=/dev/stderr
; stderr_logfile_maxbytes=0
